# Docs/Service/Handler_Parameter.service
[Unit]
Description=PostOffice Parameter Data Handler Service
After=network.target

[Service]
WorkingDirectory=/root/PostOffice
ExecStart=/usr/bin/python3 /root/PostOffice/Handler/Handler_Parameter.py
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target