# Docs/Service/Hardware.service
[Unit]
Description=Hardware API Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/PostOffice
ExecStart=gunicorn -c /root/PostOffice/Setup/Service/Hardware_config.py Hardware:Hardware

[Install]
WantedBy=multi-user.target