# Docs/Service/PostOffice.service
[Unit]
Description=PostOffice Hardware API Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/home/postoffice/PostOffice/src
ExecStart=/usr/bin/gunicorn -c /home/postoffice/PostOffice/src/Docs/Service/PostOffice_config.py PostOffice:PostOffice
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target